# deployment/compras-verdu.service
# =============================================================================
# systemd: arranque automático de Compras Verdu al iniciar el SO
# =============================================================================
# Instalación (una sola vez):
#   sudo ./deployment/instalar-servicio.sh
# O manual: editar WorkingDirectory con la ruta real del proyecto y luego:
#   sudo cp deployment/compras-verdu.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable compras-verdu.service
#
# Ver estado:  sudo systemctl status compras-verdu
# Ejecutar ahora:  sudo systemctl start compras-verdu
# =============================================================================

[Unit]
Description=Compras Verdu (Docker Compose)
Documentation=Ver INICIO.md y docs/ENV.md
After=docker.service network-online.target
Wants=docker.service network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/adm_agustin/comprasVerdu
ExecStart=/usr/bin/bash -lc 'exec ./scripts/docker-compose-env.sh up -d'
ExecStop=/usr/bin/bash -lc 'exec ./scripts/docker-compose-env.sh down'
TimeoutStartSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=compras-verdu

[Install]
WantedBy=multi-user.target
